# scripts/iplocation.conf.upstart
description "iplocation monitor"
author "wangfei@chunyu.me"

start on startup
stop on shutdown

respawn
respawn limit 2 5

console output


#setuid work
#setgid work

exec /usr/local/ip/iplocation -c /usr/local/ip/config.ini -L /usr/local/ip/log/ip.log &> /usr/local/ip/log/app.log

pre-start script
        logger "start iplocation"
end script

post-start script
        logger "iplocation start succeed"
end script

post-stop script
        logger "iplocation stop"
end script